# com.daily-paper.lint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily Paper wiki lint — runs daily at 11 PM local time.

  Installation:
    cp com.daily-paper.lint.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.daily-paper.lint.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.daily-paper.lint.plist
    rm ~/Library/LaunchAgents/com.daily-paper.lint.plist

  Logs:
    tail -f /Users/xinzhang/daily_paper/logs/lint.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daily-paper.lint</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/xinzhang/daily_paper/wiki_lint.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/xinzhang/daily_paper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/xinzhang/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/xinzhang</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/xinzhang/daily_paper/logs/lint.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/xinzhang/daily_paper/logs/lint_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
