# com.paper-search-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paper-search-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/xinzhang/daily_paper/supervisor.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/xinzhang/daily_paper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/Caskroom/miniforge/base/envs/paper_server/bin:/Users/xinzhang/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/xinzhang</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/xinzhang/daily_paper/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/xinzhang/daily_paper/logs/stderr.log</string>
</dict>
</plist>
